# documents/docs-sync.path
[Path]
PathChanged=/home/indy/Documents

[Install]
WantedBy=default.target
