# documents/docs-sync.timer
[Timer]
OnActiveSec=60s
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
